# fabdeploy/telegrambot.service
[Unit]
Description=Lexa Telegram Bot service
After=multi-user.target

[Service]
User=alexey
Group=users
WorkingDirectory=/home/alexey/projects/autopublisher
ExecStart=/home/alexey/.pyenv/versions/autopublisher/bin/python3 /home/alexey/projects/autopublisher/telegrambot.py
Type=simple
NotifyAccess=all
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
